# systemd/systemd/user/opencode.service
[Unit]
Description=OpenCode Server (fish + mise, serve)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/dotfiles/opencode/opencode/server.fish
Restart=always
RestartSec=2

Environment=PATH=/usr/local/bin:/usr/bin:/bin

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
